# systemd/litestream.service
# /etc/systemd/system/litestream.service — Story 1.c.3 (AR-014, AO-38/AO-51, D-04.21).
#
# Q-C3-1 = (b) SERVIÇO SEPARADO: o Litestream corre como daemon independente
# (`litestream replicate`), desacoplado do worker. `hdd-worker.service` declara
# Requires=/After=litestream.service → a replicação arranca ANTES do worker e o
# worker não corre sem durabilidade activa (fail-closed na durabilidade: não
# escrever state sem réplica — defesa do crash recovery do Epic 5).
#
# Difere do canon AO-51 (wrapper `litestream run -- <worker>`, 1 processo) por
# decisão consciente do operador alinhada com o files_created do StorySpec.
#
# CREDENCIAIS R2 (Q-C3-3): EnvironmentFile PRÓPRIO 0600, separado do secrets.env
# do worker (separação de concerns; o worker nunca vê creds R2). Litestream lê
# LITESTREAM_ACCESS_KEY_ID / LITESTREAM_SECRET_ACCESS_KEY automaticamente.
#
# Pré-requisito: binário litestream instalado (GitHub releases, versão pinada —
# ver runbook; Renovate nunca automerge — architecture.md:736).

[Unit]
Description=Litestream replication (HDD state -> R2 EU)
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/litestream.yml
ConditionPathExists=/etc/litestream/litestream.env

[Service]
Type=simple
User=hdd-worker
EnvironmentFile=/etc/litestream/litestream.env
# Gate de permissão das creds R2 (mesmo padrão de 1.c.2 / NFR-S1): recusa
# arrancar se o EnvironmentFile não for 0600.
ExecStartPre=/bin/bash -c 'test "$(stat -c %a /etc/litestream/litestream.env)" = "600"'
ExecStart=/usr/local/bin/litestream replicate -config /etc/litestream.yml
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
